# patch/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# We are in lanapmode, so why do we need dnsmasq?
# NOTE: This bricks the router from resolving domain name on the host, but again
#       why would an AP need to connect to the Internet?
/etc/init.d/dnsmasq stop

exit 0
